# qwen35-mtp.service
[Unit]
Description=Qwen3.6 35B-MoE Architect on 8154
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment="CUDA_SCALE_LAUNCH_QUEUES=4x"
Environment="GGML_CUDA_GRAPH_OPT=1"
Environment="GGML_CUDA_FORCE_CUBLAS_COMPUTE_16F=1"

ExecStart=/home/radek/src/llama.cpp-mtp/build/bin/llama-server \
  -hf unsloth/Qwen3.6-35B-A3B-MTP-GGUF:Q4_K_XL \
  --alias qwen3.6-35b-architect \
  --host 127.0.0.1 \
  --port 8154 \
  -ngl 999 \
  -fa on \
  -c 262144 \
  --cache-type-k q8_0 \
  --cache-type-v q8_0 \
  --kv-unified \
  --cache-ram 49152 \
  --cache-idle-slots \
  --ctx-checkpoints 64 \
  --cache-reuse 1024 \
  -b 16384 \
  -ub 4096 \
  --parallel 1 \
  --threads 8 \
  --cpu-range 0-9 \
  --threads-batch 16 \
  --cpu-range-batch 10-19 \
  --threads-http 4 \
  --prio 2 \
  --poll 50 \
  --no-mmap \
  --reasoning on \
  --reasoning-budget 4000 \
  --reasoning-format deepseek \
  --jinja \
  --temp 1.0 \
  --top-k 20 \
  --top-p 0.95 \
  --min-p 0.0 \
  --presence-penalty 0.5 \
  --keep -1 \
  --metrics \
  --slots

Restart=on-failure
RestartSec=5
LimitNOFILE=1048576

[Install]
WantedBy=default.target
