# drop-ins/qwen35-mtp.service.d/override.conf
# managed-by:harden-llm-stack
# Generated 2026-05-19T21:56:36-07:00 by harden-llm-stack.sh
# Revert: ./harden-llm-stack.sh --revert  (or rm this whole .d dir)

[Unit]
Conflicts=qwen36-27b-smart.service qwen36-35b-fast.service qwen27-mtp.service gemma-31b.service gptoss-20b.service

[Service]
# memory cap — when hit, kernel OOMs this cgroup only, host untouched
MemoryHigh=70G
MemoryMax=80G
OOMPolicy=stop
OOMScoreAdjust=200

# restart with a hard burst limit (3 fails in 10 min → systemd gives up)
Restart=on-failure
RestartSec=15
StartLimitBurst=3
StartLimitIntervalSec=600
